# conf/supervisord.conf
[supervisord]
nodaemon=true

[program:phantom_tank]
directory=/app/python
command=python3 phantom_tank.py

